# deploy/moalim-backend.service
[Unit]
Description=Moalim FastAPI backend (uvicorn)
After=network.target redis-server.service
Wants=redis-server.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/moalim/backend
EnvironmentFile=/root/moalim/backend/.env
ExecStart=/root/moalim/backend/.venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 2 \
    --proxy-headers \
    --forwarded-allow-ips="*" \
    --no-access-log
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=moalim-backend

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/root/moalim

[Install]
WantedBy=multi-user.target
