# deploy/moalim-backend.service
[Unit]
Description=Moalim FastAPI backend (uvicorn) — isolated service
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/moalim/backend
EnvironmentFile=/root/moalim/backend/.env
ExecStart=/root/moalim/backend/.venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 1 \
    --proxy-headers \
    --forwarded-allow-ips="*" \
    --no-access-log
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=moalim-backend

# Sécurité — n'affecte que ce service
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
